# scripts/cron.d/wb-refresh-tariffs
# Авто-обновление тарифов складов WB (КС) для юнит-экономики.
# Аналог /etc/cron.d/ue-sanity-check — работает НАПРЯМУЮ, не зависит от Laravel
# schedule:run (который на этом сервере может быть не настроен).
#
# Установка на сервере (один раз):
#   sudo cp /var/www/products-backend/scripts/cron.d/wb-refresh-tariffs /etc/cron.d/wb-refresh-tariffs
#   sudo chown root:root /etc/cron.d/wb-refresh-tariffs
#   sudo chmod 644 /etc/cron.d/wb-refresh-tariffs
# cron.d подхватывается автоматически, перезапуск не нужен. Проверка:
#   sudo run-parts --report /etc/cron.d  ||  grep wb-refresh /var/log/syslog
#
# ВАЖНО: при использовании этого cron.d-файла отключите дублирующее расписание
# Laravel, чтобы команда не запускалась дважды (она идемпотентна, но лишние
# запросы к WB API не нужны): в .env поставьте  WB_TARIFF_SCHEDULE=false
# и выполните  php artisan config:cache .
#
# При необходимости поправьте под своё окружение: пользователь (www-data),
# путь к php (which php), путь к проекту.

SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
MAILTO=""

# Ежедневно в 05:30 МСК — обновить КС всех активных WB-интеграций
# (box-снапшоты + inventory_warehouses) и пересобрать кэш юнит-экономики.
30 5 * * * www-data cd /var/www/products-backend && /usr/bin/php artisan wb:refresh-tariffs >> /var/www/products-backend/storage/logs/wb-tariff-refresh.log 2>&1
